# example.service
[Unit]
Description=Service Description
After=network.target

[Service]
User=usename
WorkingDirectory=/path_to_app
ExecStart=/path_to_app/venv/bin/python -m uvicorn api:app --host 0.0.0.0 --port 8585
Restart=always

[Install]
WantedBy=multi-user.target